# provision_rpi/templates/pigpiod.service.j2
[Unit]
Description=Daemon required to control GPIO pins via pigpio
After=network.target

[Service]
ExecStart=/usr/local/bin/pigpiod -l
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -SIGINT $MAINPID
Restart=on-failure
Type=forking

# This prevents systemd from sending SIGCONT (signal 18)
KillMode=process
SendSIGKILL=no

[Install]
WantedBy=multi-user.target
